# webapp/deploy/systemd/dailycup-analytics-healthcheck.timer
[Unit]
Description=DailyCup Analytics Refresh Healthcheck hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target